# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

# Optional: Redis (only needed if SMART_KANBAN_MODE=production)
# In local mode, the app uses in-memory fallback
[program:redis]
command=redis-server --bind 0.0.0.0 --port 6379
autostart=false
autorestart=true
stdout_logfile=/workspace/data/logs/redis.log
stderr_logfile=/workspace/data/logs/redis.err
priority=1

# Backend API (FastAPI with uvicorn)
[program:backend]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000 --reload
directory=/workspace/backend
autostart=true
autorestart=true
stdout_logfile=/workspace/data/logs/backend.log
stderr_logfile=/workspace/data/logs/backend.err
priority=10
environment=PYTHONUNBUFFERED="1",SMART_KANBAN_MODE="%(ENV_SMART_KANBAN_MODE)s",DATABASE_URL="%(ENV_DATABASE_URL)s",FRONTEND_URL="%(ENV_FRONTEND_URL)s"

# Celery Worker
[program:worker]
command=celery -A app.celery_app worker --loglevel=info --pool=solo
directory=/workspace/backend
autostart=true
autorestart=true
stdout_logfile=/workspace/data/logs/worker.log
stderr_logfile=/workspace/data/logs/worker.err
priority=20
environment=PYTHONUNBUFFERED="1",SMART_KANBAN_MODE="%(ENV_SMART_KANBAN_MODE)s",DATABASE_URL="%(ENV_DATABASE_URL)s"

# Frontend (serve built files with a simple HTTP server)
[program:frontend]
command=npx serve -s dist -l 3000
directory=/workspace/frontend
autostart=true
autorestart=true
stdout_logfile=/workspace/data/logs/frontend.log
stderr_logfile=/workspace/data/logs/frontend.err
priority=30
